# modules/minecraft_serv/temlates/minecraft.service.epp
[Unit]
Description=minecraft_server

[Service]
ExecStart=/usr/bin/java -Xmx<%= $xmx %> -Xms<%= $xms %> -jar /opt/minecraft/server.jar nogui
ExecStop=/usr/bin/java -Xmx<%= $xmx %> -Xms<%= $xms %> -jar /opt/minecraft/server.jar stop
[Install]
WantedBy=multi-user.target
